# deploy/ansible/roles/web/templates/babylon-web.service.j2
{# spec 061 T021: systemd unit for the Django/Gunicorn worker.
   Companion to (and intended replacement for) the supervisor-managed
   gunicorn_start.j2 wrapper. Per Constitution X.4, systemd is the sole
   supervisor on the target host; the supervisor configuration in
   roles/web/tasks/setup_supervisor.yml is pre-existing tech debt
   flagged in the spec 061 plan.md as out of scope for this feature.

   Exponential restart backoff requires systemd 254+ (Debian 13 /
   Ubuntu 24.04+). Verified on dev host (systemd 257) per T005.

   Per research.md R3 + R4: do NOT use --preload-app — fork-safety
   collides with psycopg ConnectionPool file descriptors. The in-process
   3-attempt retry inside babylon.web.game.apps.GameConfig handles
   transient sub-second blips; this systemd unit's exponential
   RestartSec backoff handles longer-timescale recovery. #}
[Unit]
Description=Babylon Gunicorn application server (Django + EngineBridge)
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service
Documentation=file:///opt/babylon/specs/061-real-backend-wireup/research.md
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=notify
User={{ gunicorn_user }}
Group={{ gunicorn_group }}
WorkingDirectory={{ project_path }}/web
Environment="PYTHONPATH={{ project_path }}:{{ project_path }}/src"
Environment="DJANGO_SETTINGS_MODULE=babylon_web.settings.production"
EnvironmentFile=-{{ virtualenv_path }}/etc/babylon-web.env

ExecStart={{ virtualenv_path }}/bin/gunicorn \
    --name {{ application_name }} \
    --workers {{ gunicorn_num_workers }} \
    --max-requests {{ gunicorn_max_requests }} \
    --timeout {{ gunicorn_timeout_seconds|default(30) }} \
    --bind unix:{{ virtualenv_path }}/run/gunicorn.sock \
    --log-level info \
    {{ application_name }}.wsgi

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30
PrivateTmp=true

# Spec 061 FR-007: hybrid retry-then-exit. The Django ready() hook does
# 3 fast in-process retries (1s, 2s) for transient DB blips, then exits
# non-zero. systemd handles persistent failures with exponential
# RestartSec → RestartMaxDelaySec.
Restart=on-failure
RestartSec=5s
RestartSteps=4
RestartMaxDelaySec=60s

[Install]
WantedBy=multi-user.target
